# setup/crontab
* * * * * /home/pi/timelapse/scripts/verify.sh
@reboot /home/pi/timelapse/scripts/verify.sh
